# app/s6-rc.d/sshd-prepare/up
#!/command/execlineb -P

with-contenv

importas ssh_public_key SSH_PUBLIC_KEY

if { mkdir -p -m 750 /var/run/sshd }
if { mkdir -p -m 700 /root/.ssh }
if { redirfd -w 1 /root/.ssh/authorized_keys echo "${ssh_public_key}" }
if { chmod 400 /root/.ssh/authorized_keys }
ssh-keygen -A 


